# varnish-elb.service
[Unit]
Description=Generate Varnish backend director for ELB

[Service]
Type=oneshot
ExecStart=/usr/local/bin/generate_backends
