# varnish-reload.service
[Unit]
Requires=varnish.service
Description=Simply reloads Varnish

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reload varnish
